# download/etc/cron.d/hermes_quarantined_report_8_hours_job
15 */8 * * * root /usr/bin/curl --silent http://localhost:8888/schedule/quarantine_report.cfm?frequency=8 &>/dev/null

