# download/etc/cron.d/hermes_message_cleanup_job
30 1 * * * root /usr/bin/curl --silent http://localhost:8888/schedule/message_cleanup.cfm &>/dev/null

