# download/etc/cron.d/hermes_quarantined_report_2_hours_job
15 */2 * * * root /usr/bin/curl --silent http://localhost:8888/schedule/quarantine_report.cfm?frequency=2 &>/dev/null

