# download/etc/cron.d/hermes_quarantined_report_daily_job
30 0 * * * root /usr/bin/curl --silent http://localhost:8888/schedule/quarantine_report_daily.cfm &>/dev/null

