# download/etc/cron.d/hermes_quarantined_report_4_hours_job
15 */4 * * * root /usr/bin/curl --silent http://localhost:8888/schedule/quarantine_report.cfm?frequency=4 &>/dev/null

